# roles/tilemaker/templates/tilemaker.timer
[Unit]
Description=Rebuild mbtiles every night

[Timer]
OnCalendar=*-*-* 04:15:00
Persistent=true

[Install]
WantedBy=timers.target
