# roles/digitransit/templates/digitransit-restart.timer
[Unit]
Description=Restart container every night

[Timer]
OnCalendar=*-*-* 04:45:00
Persistent=true

[Install]
WantedBy=timers.target
