# roles/digitransit/templates/digitransit-restart.timer
[Unit]
Description=Restart container every night

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
