# roles/data-builder/templates/data-builder.timer
[Unit]
Description=Build a container every night

[Timer]
# run it every night at 11
OnCalendar=*-*-* 11:00:00

[Install]
WantedBy=timers.target
