# roles/services/templates/systemd/layer-categories-lbg.timer
[Unit]
Description=Run my service periodically

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target