# roles/tileserver/templates/tileserver-restart.timer
[Unit]
Description=Restart container every night

[Timer]
OnCalendar=*-*-* 05:00:00
Persistent=true

[Install]
WantedBy=timers.target
