# roles/poi/templates/systemd/import-pois.timer
[Unit]
Description=Timer for poi service

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target