# roles/tileserver/templates/systemd/tileserver-restart.timer
[Unit]
Description=Restart container every night

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true

[Install]
WantedBy=timers.target
