# roles/thingsboard-enhancer/templates/thingsboard-enhancer.timer
[Unit]
Description=Run thingsboard-enhancer.service every minute

[Timer]
# every 15 seconds
#OnCalendar=*:*:0/15
#
# every minute
OnCalendar=*:0/1
Persistent=true

[Install]
WantedBy=timers.target
