# roles/delay-prediction-service/templates/delay-prediction-cleanup.timer
[Unit]
Description=delete old delay-prediction-service data every night

[Timer]
OnCalendar=*-*-* 04:30:00
Persistent=true

[Install]
WantedBy=timers.target
