# roles/delay-prediction-service/templates/delay-prediction.service
[Unit]
Description=pull bus positions from Thingsboard, predict delays & positions, send GTFS-RT to MQTT
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Restart=always
RestartSec=15

ExecStart=/var/delay-prediction/delay-prediction-service delay-prediction-service

[Install]
WantedBy=multi-user.target
