# roles/digitransit-next/templates/digitransit-restart.timer
[Unit]
Description=Restart digitransit every night

[Timer]
OnCalendar=*-*-* 02:45:00
Persistent=true

[Install]
WantedBy=timers.target
