# roles/photon/templates/nominatim-update.timer
[Unit]
Description=Run nominatim-update daily

[Timer]
OnCalendar=*-*-* 2:00:00
Persistent=true

[Install]
WantedBy=timers.target
