# roles/photon/templates/nominatim-update.service
[Unit]
Description=Update nominatim with the update files from Geofabrik
Requires=docker.service
After=docker.service

[Service]
WorkingDirectory=/var/nominatim/
ExecStart=/var/nominatim/update-nominatim

#ExecStopPost=-/usr/local/bin/send-to-matrix "❌ Nominatim update failed on {{ inventory_hostname }}"

[Install]
WantedBy=multi-user.target
