# roles/photon/templates/photon-update.timer
[Unit]
Description=Run nominatim-update daily

[Timer]
OnCalendar=*-*-* 5:15:00
Persistent=true

[Install]
WantedBy=timers.target
