# roles/digitransit/templates/systemd/digitransit-restart.timer
[Unit]
Description=Restart digitransit every night

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
