# roles/base/templates/docker-prune.timer
[Unit]
Description=Run docker prune daily

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
